# rasbuild/files/etc/systemd/system/initialsetup.service
[Unit]
Description=Raspberry Server Initialize Service

[Service]
ExecStart=/bin/bash /root/startup.sh

[Install]
WantedBy=multi-user.target
